# systemd/paribus-bulk.service
[Unit]
Description=Paribus Hospital Bulk Processor
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/paribus-hospital-directory
EnvironmentFile=-/home/ubuntu/paribus-hospital-directory/.env
User=ubuntu
Group=ubuntu
# Explicitly expose production environment to the process
Environment=APP_ENV=production
# Use Gunicorn with Uvicorn workers for a production-grade process manager.
# Bind to the port expected by the upstream nginx configuration.
ExecStart=/usr/bin/env poetry run gunicorn -k uvicorn.workers.UvicornWorker "app.main:app" \
	--bind 0.0.0.0:7103 \
	--workers 4 \
	--log-level info \
	--access-logfile - \
	--error-logfile -
Restart=on-failure
RestartSec=5
# Make sure child processes are killed on stop and allow graceful shutdown
KillMode=process
TimeoutStopSec=30s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
